# scripts/kotekan.service
[Unit]
Description=kotekan
After=syslog.target network.target

[Service]
Type=simple
Environment="DISPLAY=:0"
WorkingDirectory=/var/lib/kotekan
ExecStart=/usr/sbin/kotekan -c /etc/kotekan/kotekan_opencl_(?(%H==0[1-9]|10)280x|fury).yaml
PIDFile=/var/run/kotekan.pid

[Install]
WantedBy=default.target
